# config/systemd/user/backblaze.timer
[Unit]
Description=Backup with Backblaze once per day

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target
