# config/systemd/user/gcal.timer
[Unit]
Description=Google Calendar Agenda Sync

[Timer]
OnCalendar=*:0/5

[Install]
WantedBy=timers.target
