# config/systemd/user/newsboat.timer
[Unit]
Description=Refresh Newsboat

[Timer]
OnCalendar=hourly

[Install]
WantedBy=timers.target
